# bin/server/daemon
#!/bin/sh

NODE_ENV="production"
NODE_NAME="dportal"
NODE_APP="js/serv.js"
APP_DIR="/home/ctrack/D-Portal/$NODE_NAME"
APP_URL="http://localhost:1408/"

PID_FILE=$APP_DIR/$NODE_ENV/$NODE_NAME.pid
LOG_FILE=$APP_DIR/$NODE_ENV/$NODE_NAME.log
CRONLOG_FILE=$APP_DIR/$NODE_ENV/$NODE_NAME.cron.log

CRONLOG_EMAILS="krissd@gmail.com,notshi@gmail.com"

NODE_EXEC=`which node`

###############

# REDHAT chkconfig header

# chkconfig: - 58 74
# description: this is the script for starting a node app on boot and keeping it running via cron.
### BEGIN INIT INFO
# Provides: node
# Required-Start:    $network $remote_fs $local_fs
# Required-Stop:     $network $remote_fs $local_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: start and stop node
# Description: Node process for app
### END INIT INFO

start_app (){
	if [ -f $PID_FILE ]
	then
		PID=`cat $PID_FILE`
		if [ -z "`ps ef -A | awk '{print $1}' | grep "^$PID$"`" ]
		then
			echo "$PID_FILE exists, and process is already running."
			return
		else
			echo "$PID_FILE exists, but process is not running"
			rm -f $PID_FILE;
 		fi
	fi
	echo "Starting node app $NODE_NAME..."
	cd $APP_DIR
	$NODE_EXEC ./$NODE_APP  1>$LOG_FILE 2>&1 &
	echo $! > $PID_FILE;
}

stop_app (){
	if [ ! -f $PID_FILE ]
	then
		echo "$PID_FILE does not exist, process is not running"
	else
		echo "Stopping $APP_DIR/$NODE_APP ..."
		echo "Killing `cat $PID_FILE`"
		kill `cat $PID_FILE`;
		rm -f $PID_FILE;
		echo "Node stopped"
	fi
}

status_app (){
	if [ -f $PID_FILE ]
	then
		PID=`cat $PID_FILE`
		if [ -z "`ps ef -A | awk '{print $1}' | grep "^$PID$"`" ]
		then
			echo "Node app $NODE_NAME stopped but pid file exists"
		else
			echo "Node app $NODE_NAME running with pid $PID"
		fi
	else
		echo "Node app $NODE_NAME stopped"
	fi
}

check_app (){
	if [ -f $PID_FILE ]
	then
		PID=`cat $PID_FILE`
		if [ -z "`ps ef -A | awk '{print $1}' | grep "^$PID$"`" ]
		then
			echo "Node app $NODE_NAME stopped but pid file exists"
			stop_app
			start_app
			status_app
			exit 10
		else
			echo "Node app $NODE_NAME running with pid $PID"
			if [ "`curl -s -o /dev/null -w "%{http_code}" $APP_URL`" = "200" ]
			then
				echo "Node app $NODE_NAME responding OK"
				exit 0
			else
				echo "Node app $NODE_NAME not responding"
				stop_app
				start_app
				status_app
				exit 20
			fi
		fi
	else
		echo "Node app $NODE_NAME stopped"
		start_app
		status_app
		exit 10
	fi
}

case "$1" in
	start)
		start_app
		status_app
	;;

	stop)
		stop_app
		status_app
	;;

	restart)
		stop_app
		start_app
		status_app
	;;

	status)
		status_app
	;;

	check)
		check_app
	;;

	cron)
		$0 check >$CRONLOG_FILE
		if [[ $? -ne 0 ]]
		then
			mail -s"$NODE_NAME on $HOSTNAME restarting" $CRONLOG_EMAILS <$CRONLOG_FILE
		fi
	;;

	*)
		echo "Usage: $0 {start|stop|restart|status|check|cron}"
		exit 1
	;;
esac
